# crontab
*/30 * * * * cd /data && python3 /app/main.py
